# local_root/_rootfs_1.extracted/squashfs-root/usr/lib/opkg/info/base-files.list
/usr/lib/common.awk
/usr/share/udhcpc/default.script
/var
/init
/sbin/wifi
/sbin/routerCfg.sh
/sbin/firstboot
/sbin/hotplug-call
/sbin/led.sh
/sbin/ifup
/sbin/saveCurrentChannel
/sbin/ifdown
/sbin/defaultCfg.sh
/sbin/flashFW.sh
/sbin/iperfServer.sh
/sbin/mount_root
/sbin/sysupgrade
/bin/login.sh
/bin/endpoint
/bin/ipcalc.sh
/bin/start_endpoint
/bin/endpoint.ini
/bin/echr.msg
/rom/note
/lib/wifi/ralink.sh
/lib/network/config.sh
/lib/firstboot/40_no_fo_copy_ramoverlay
/lib/firstboot/30_reset_copy_rom
/lib/firstboot/40_copy_ramoverlay
/lib/firstboot/30_no_fo_pivot
/lib/firstboot/20_reset_clear_jffs
/lib/firstboot/10_no_fo_clear_overlay
/lib/firstboot/10_determine_parts
/lib/firstboot/05_firstboot_skip
/lib/firstboot/20_has_mini_fo
/lib/firstboot/50_pivot
/lib/firstboot/10_reset_has_mini_fo
/lib/firstboot/99_10_with_fo_cleanup
/lib/firstboot/30_is_rootfs_mounted
/lib/firstboot/20_no_fo_mount_jffs
/lib/firstboot/99_10_no_fo_cleanup
/lib/functions/boot.sh
/lib/preinit/90_mount_no_jffs2
/lib/preinit/20_device_fs_mount
/lib/preinit/90_init_console
/lib/preinit/90_restore_config
/lib/preinit/05_mount_skip
/lib/preinit/10_essential_fs
/lib/preinit/10_check_for_mtd
/lib/preinit/10_indicate_preinit
/lib/preinit/30_device_fs_daemons
/lib/preinit/40_mount_devpts
/lib/preinit/50_indicate_regular_preinit
/lib/preinit/60_init_hotplug
/lib/preinit/70_pivot_jffs2_root
/lib/preinit/50_choose_console
/lib/preinit/70_initramfs_test
/lib/preinit/40_mount_jffs2
/lib/preinit/40_init_shm
/lib/preinit/20_check_jffs2_ready
/lib/preinit/99_10_mount_no_mtd
/lib/preinit/99_10_run_init
/lib/preinit/80_mount_root
/lib/upgrade/common.sh
/etc/config/wireless
/etc/config/network
/etc/config/system
/etc/config/dhcp
/etc/proc_sys_kernel_hotplug.sh
/etc/inittab
/etc/resolv.conf
/etc/TZ
/etc/functions.sh
/etc/protocols
/etc/passwd
/etc/fstab
/etc/check_mount_media.sh
/etc/em357-fw/em357-nodetest-use-with-bootloader.ebl
/etc/em357-fw/em357-ncp-uart-rts-cts-use-with-serial-uart-bl-511.ebl
/etc/rc.local
/etc/force_mount_media.sh
/etc/openwrt_version
/etc/rc.common
/etc/preinit
/etc/hosts
/etc/hotplug.d/ieee1394/10-ieee1394
/etc/hotplug.d/usb/10-usb
/etc/hotplug.d/iface/10-routes
/etc/hotplug.d/iface/00-netstate
/etc/hotplug.d/iface/40-rdate
/etc/hotplug.d/net/10-net
/etc/hotplug2-common.rules
/etc/diag.sh
/etc/mdev_mount_umount.sh
/etc/banner
/etc/shells
/etc/sysctl.conf
/etc/openwrt_release
/etc/mtab
/etc/init.d/network
/etc/init.d/rcS
/etc/init.d/done
/etc/init.d/usb
/etc/init.d/umount
/etc/init.d/boot
/etc/init.d/watchdog
/etc/init.d/sysctl
/etc/init.d/led
/etc/group
/etc/Wireless/RT2860/RT2860.dat
/etc/Wireless/RT2860/RT2860.dat.base
/etc/profile
/etc/ntp.sh
/etc/hotplug2-init.rules
