# local_root/_rootfs_1.extracted/squashfs-root/usr/lib/opkg/info/busybox.list
/usr/sbin/crond
/usr/sbin/arping
/usr/sbin/chroot
/usr/sbin/tftpd
/usr/sbin/telnetd
/usr/sbin/brctl
/usr/sbin/ntpd
/usr/sbin/killall5
/usr/sbin/rdate
/usr/bin/nslookup
/usr/bin/wget
/usr/bin/tail
/usr/bin/uptime
/usr/bin/awk
/usr/bin/wc
/usr/bin/tftp
/usr/bin/free
/usr/bin/bunzip2
/usr/bin/cmp
/usr/bin/tr
/usr/bin/clear
/usr/bin/time
/usr/bin/passwd
/usr/bin/uniq
/usr/bin/diff
/usr/bin/sort
/usr/bin/md5sum
/usr/bin/pgrep
/usr/bin/xargs
/usr/bin/strings
/usr/bin/[[
/usr/bin/less
/usr/bin/tee
/usr/bin/readlink
/usr/bin/top
/usr/bin/basename
/usr/bin/id
/usr/bin/nc
/usr/bin/pkill
/usr/bin/crontab
/usr/bin/cksum
/usr/bin/yes
/usr/bin/printf
/usr/bin/expr
/usr/bin/[
/usr/bin/du
/usr/bin/mkfifo
/usr/bin/hexdump
/usr/bin/seq
/usr/bin/test
/usr/bin/reset
/usr/bin/dirname
/usr/bin/traceroute
/usr/bin/cut
/usr/bin/mesg
/usr/bin/which
/usr/bin/telnet
/usr/bin/killall
/usr/bin/env
/usr/bin/hostid
/usr/bin/find
/usr/bin/bzcat
/usr/bin/logger
/usr/bin/head
/sbin/halt
/sbin/init
/sbin/blkid
/sbin/insmod
/sbin/vconfig
/sbin/mdev
/sbin/getty
/sbin/mkswap
/sbin/hwclock
/sbin/devmem
/sbin/udhcpc
/sbin/pivot_root
/sbin/poweroff
/sbin/ifconfig
/sbin/start-stop-daemon
/sbin/swapoff
/sbin/klogd
/sbin/lsmod
/sbin/reboot
/sbin/syslogd
/sbin/fdisk
/sbin/rmmod
/sbin/logread
/sbin/watchdog
/sbin/sysctl
/sbin/swapon
/sbin/route
/sbin/switch_root
/bin/chgrp
/bin/uname
/bin/sync
/bin/netmsg
/bin/dmesg
/bin/login
/bin/mknod
/bin/mkdir
/bin/kill
/bin/grep
/bin/tar
/bin/df
/bin/zcat
/bin/ls
/bin/fgrep
/bin/sh
/bin/pidof
/bin/ping
/bin/egrep
/bin/touch
/bin/cat
/bin/chmod
/bin/ln
/bin/busybox
/bin/ash
/bin/pwd
/bin/gzip
/bin/date
/bin/true
/bin/cp
/bin/nice
/bin/echo
/bin/fsync
/bin/vi
/bin/umount
/bin/netstat
/bin/sed
/bin/mv
/bin/dd
/bin/sleep
/bin/rm
/bin/mktemp
/bin/chown
/bin/false
/bin/mount
/bin/gunzip
/bin/rmdir
/bin/ping6
/bin/ps
/bin/lock
/etc/init.d/cron
/etc/init.d/sysntpd
/etc/init.d/telnet
