# local_root/_rootfs_1.extracted/squashfs-root/usr/lib/opkg/info/WeMo.list
/sbin/startWemo.sh
/sbin/wemoApp
/etc/ver.txt
/etc/sensor.jpg
/etc/icon.jpg
/etc/init.d/wemo
